# deployer/cron.d
MAILTO=andychilton@gmail.com

# mins hour dofm mnth dofw   user          command
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
     0    *    *    *    *   chilts        find /var/cache/javascript-minifier.com/ -name '*.js' -mmin +1440 -delete
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
